# ansible/roles/create_service/templates/fibonacci.service.j2
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/node /home/vagrant/TP-DEVOPS start
ExecStop=/usr/bin/node /home/vagrant/TP-DEVOPS stop
